# weather-clock.service
[Unit]
Description = Weather Clock

[Service]
ExecStart = /usr/local/bin/weather-clock
Restart = never

[Install]
WantedBy = default.target
